# cron/sync_users.cron
*/5 * * * * root ( cd /usr/local/ciconnect && ./sync_users.sh -u root.atlas-af -g root.atlas-af -e https://api.ci-connect.net:18080 ) >> /var/log/provisioner.log 2>&1
